# etc/systemd/user/pulseaudio.service.d/30-wait.conf
[Service]
# Sleep for 2 seconds before restarting the service
RestartSec=2
